# setup/certbot/certbot.service
[Unit]
Description=Certbot TLS Renewal
Documentation=https://certbot.eff.org/docs

[Service]
Type=oneshot
ExecStart=sudo bash /foundrycron/certbot.sh
