# setup/certbot/certbot_start.timer
[Unit]
Description=Run certbot once at start without a random delay

[Timer]
Unit=certbot.service
OnActiveSec=120s

[Install]
WantedBy=timers.target
