# setup/aws/dynamic_dns.service
[Unit]
Description=Set dynamic DNS in case it changes

[Service]
ExecStart=sudo bash /foundrycron/dynamic_dns.sh
