# setup/aws/hosted_zone_ip.service
[Unit]
Description=Set dynamic DNS in case it changes

[Service]
Type=oneshot
ExecStart=sudo bash /foundrycron/hosted_zone_ip.sh
