# files/certbot/certbot_start.timer
[Unit]
Description=Run certbot without a random delay

[Timer]
OnActiveSec=30s

[Install]
WantedBy=timers.target