# files/certbot/certbot_start.service
[Unit]
Description=Certbot SSL
Documentation=https://certbot.eff.org/docs

[Service]
Type=oneshot
ExecStart=sudo bash /foundrycron/certbot.sh