# files/certbot/certbot_renew.service
[Unit]
Description=Certbot
Documentation=https://certbot.eff.org/docs

[Service]
Type=oneshot
ExecStart=sudo bash /foundrycron/certbot_renew.sh